# com.clipsautomation.tiktok.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.clipsautomation.tiktok</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/gustavorocha/Projetos/ClipsAutomation/venv/bin/python</string>
        <string>/Users/gustavorocha/Projetos/ClipsAutomation/tiktok_runner.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/gustavorocha/Projetos/ClipsAutomation</string>

    <!-- Roda a cada 2 horas -->
    <key>StartInterval</key>
    <integer>7200</integer>

    <key>StandardOutPath</key>
    <string>/Users/gustavorocha/Projetos/ClipsAutomation/logs/tiktok.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/gustavorocha/Projetos/ClipsAutomation/logs/tiktok.log</string>

    <!-- Não reinicia automaticamente em caso de falha -->
    <key>KeepAlive</key>
    <false/>

    <!-- Só roda se o Mac estiver ligado; recupera execuções perdidas -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
